# build/com.wattsadvisor.muskmarket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wattsadvisor.muskmarket</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/kiran/Desktop/Company Dossiers /Elon Musk/build/refresh_market_site.py</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MUSK_WATCH_NTFY</key>
        <string>https://ntfy.sh/musk-sec-9053334806</string>
        <key>MUSK_WATCH_GIT_PUSH</key>
        <string>1</string>
        <key>MUSK_LINKEDIN_ALERT</key>
        <string>1</string>
        <key>MUSK_DOSSIER_URL</key>
        <string>https://kaydeep0.github.io/musk-empire-dossier/</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kiran/Desktop/Company Dossiers /Elon Musk/data/market.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kiran/Desktop/Company Dossiers /Elon Musk/data/market.err.log</string>
</dict>
</plist>
